# init-scripts/facebook-scraper.conf
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/xvfb.log
stderr_logfile=/var/log/facebook-scraper/xvfb_err.log
priority=10

[program:celery_beat]
command=/opt/facebook-scraper/venv/bin/celery -A multi_queue_worker beat --loglevel=info
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/celery_beat.log
stderr_logfile=/var/log/facebook-scraper/celery_beat_err.log
environment=PYTHONPATH="/opt/facebook-scraper",DISPLAY=":99"

[program:celery_worker_scan]
command=/opt/facebook-scraper/init-scripts/celery_worker_dynamic.sh scan
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/worker_scan.log
stderr_logfile=/var/log/facebook-scraper/worker_scan_err.log
environment=PYTHONPATH="/opt/facebook-scraper",DB_HOST="localhost",DB_PORT="5432",DB_USER="fb_user",DB_PASSWORD="fbscraper_e276a5b73b04742c",DB_NAME="facebook_scraper",REDIS_HOST="localhost",REDIS_PORT="6379",REDIS_DB="0",DOCKER_ENV="false",API_BASE_URL="http://localhost:8000",DISPLAY=":99"

[program:celery_worker_maintenance]
command=/opt/facebook-scraper/init-scripts/celery_worker_dynamic.sh maintenance
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/worker_maintenance.log
stderr_logfile=/var/log/facebook-scraper/worker_maintenance_err.log
environment=PYTHONPATH="/opt/facebook-scraper",DB_HOST="localhost",DB_PORT="5432",DB_USER="fb_user",DB_PASSWORD="fbscraper_e276a5b73b04742c",DB_NAME="facebook_scraper",REDIS_HOST="localhost",REDIS_PORT="6379",REDIS_DB="0",DOCKER_ENV="false",API_BASE_URL="http://localhost:8000",DISPLAY=":99"

[program:api]
command=/opt/facebook-scraper/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/api.log
stderr_logfile=/var/log/facebook-scraper/api_err.log
environment=PYTHONPATH="/opt/facebook-scraper"

[program:streamlit]
command=/opt/facebook-scraper/venv/bin/streamlit run webapp_streamlit/app.py --server.port 8501 --server.address 0.0.0.0
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/streamlit.log
stderr_logfile=/var/log/facebook-scraper/streamlit_err.log
environment=PYTHONPATH="/opt/facebook-scraper"

[program:flower]
command=/opt/facebook-scraper/venv/bin/celery -A multi_queue_worker flower --port=5555
directory=/opt/facebook-scraper
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/flower.log
stderr_logfile=/var/log/facebook-scraper/flower_err.log
environment=PYTHONPATH="/opt/facebook-scraper"

[program:ngrok]
command=/usr/local/bin/ngrok http 8501 --config /root/.config/ngrok/ngrok.yml --log stdout
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/facebook-scraper/ngrok.log
stderr_logfile=/var/log/facebook-scraper/ngrok_err.log

[group:facebook-scraper]
programs=xvfb,celery_beat,celery_worker_scan,celery_worker_maintenance,api,streamlit,flower,ngrok

